# config/helixscreen.service
# SPDX-License-Identifier: GPL-3.0-or-later
#
# systemd service file for HelixScreen
#
# This service launches HelixScreen with an early splash screen for
# instant visual feedback during boot. The launcher starts a lightweight
# splash binary in parallel with the main app. When the main app is ready,
# it sends SIGUSR1 to the splash process, which then exits gracefully.
#
# Installation:
#   sudo cp helixscreen.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable helixscreen
#   sudo systemctl start helixscreen
#
# Logs (application logs to systemd journal automatically):
#   sudo journalctl -u helixscreen -f
#   sudo journalctl -t helix -f          # Filter by syslog identifier

[Unit]
Description=HelixScreen 3D Printer Touch UI
Documentation=https://github.com/prestonbrown/helixscreen
After=network.target

# Wait for display hardware to be ready
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=simple

# Run as the Klipper ecosystem user (templated by installer)
User=@@HELIX_USER@@
Group=@@HELIX_GROUP@@

# Working directory for config/assets
WorkingDirectory=@@INSTALL_DIR@@

# Unbind framebuffer console before starting (requires root, runs before User= drop)
# vtcon1 is the fbcon driver; vtcon0 is the dummy. Unbinding prevents kernel
# messages (dmesg, undervoltage) from painting over the UI via fbcon.
ExecStartPre=+/bin/sh -c 'for f in /sys/class/vtconsole/vtcon*/bind; do echo 0 > "$$f" 2>/dev/null || true; done'

# Use launcher script for splash + main app
ExecStart=@@INSTALL_DIR@@/bin/helix-launcher.sh

# Environment variables
# Use fbdev display backend on Pi (DRM atomic modesetting has compatibility issues
# with some SoCs like Allwinner H616). fbdev works reliably across all Pi hardware.
Environment="HELIX_DISPLAY_BACKEND=fbdev"
# Override asset directory if ui_xml/assets/config are not under WorkingDirectory
# Environment="HELIX_DATA_DIR=/usr/share/helixscreen"
# Logging: auto-detects journal on systemd systems (can override with HELIX_LOG_DEST)
# Environment="HELIX_LOG_DEST=journal"

# Restart policy - always restart unless explicitly stopped
Restart=always
RestartSec=3

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening (can be relaxed if needed)
ProtectSystem=strict
ReadWritePaths=@@INSTALL_DIR@@/config /var/log
PrivateTmp=true
NoNewPrivileges=true

# Allow framebuffer/DRM access
SupplementaryGroups=video input render tty

# Console suppression: KDSETMODE KD_GRAPHICS prevents kernel messages from
# bleeding through LVGL's partial framebuffer repaints. Requires CAP_SYS_TTY_CONFIG.
AmbientCapabilities=CAP_SYS_TTY_CONFIG

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helixscreen

[Install]
WantedBy=multi-user.target
